# templates/default/vault.conf.erb
#Upstart for Vault
description "Vault Up Start Service"

start on startup
stop on runlevel [!2345]

respawn
respawn limit 10 20


exec su - <%= node[:vault][:user] %> -c <%= node[:vault][:install_dir] %>/vault/vault -config=<%= node.vault[:conf_dir] %>/config.hcl 2>&1 >> <%= node[:vault][:logs_dir] %>/vault.log

post-stop script
  sleep 2
  if /usr/bin/pgrep -u vault; then
    /usr/bin/pgrep -u vault | /usr/bin/xargs kill -9
  fi
end script